# CMSSW_12_4_X/el8_amd64_gcc10/tools/tauolapp/1.1.8-5c4faaa3d5ddc2d72dd03a1f4952c54d/etc/profile.d/init.sh
if [ -f /work/build/CMSSW_12_4_X/el8_amd64_gcc10/tmp/el8_amd64_gcc10/external/tauolapp/1.1.8-5c4faaa3d5ddc2d72dd03a1f4952c54d/etc/profile.d/dependencies-setup.sh ]; then . /work/build/CMSSW_12_4_X/el8_amd64_gcc10/tmp/el8_amd64_gcc10/external/tauolapp/1.1.8-5c4faaa3d5ddc2d72dd03a1f4952c54d/etc/profile.d/dependencies-setup.sh; fi
TAUOLAPP_ROOT="/work/build/CMSSW_12_4_X/el8_amd64_gcc10/tmp/el8_amd64_gcc10/external/tauolapp/1.1.8-5c4faaa3d5ddc2d72dd03a1f4952c54d"
TAUOLAPP_VERSION="1.1.8-5c4faaa3d5ddc2d72dd03a1f4952c54d"
TAUOLAPP_REVISION="1"
TAUOLAPP_CATEGORY="external"
[ ! -d /work/build/CMSSW_12_4_X/el8_amd64_gcc10/tmp/el8_amd64_gcc10/external/tauolapp/1.1.8-5c4faaa3d5ddc2d72dd03a1f4952c54d/bin ] || export PATH="/work/build/CMSSW_12_4_X/el8_amd64_gcc10/tmp/el8_amd64_gcc10/external/tauolapp/1.1.8-5c4faaa3d5ddc2d72dd03a1f4952c54d/bin${PATH:+:$PATH}";
[ ! -d /work/build/CMSSW_12_4_X/el8_amd64_gcc10/tmp/el8_amd64_gcc10/external/tauolapp/1.1.8-5c4faaa3d5ddc2d72dd03a1f4952c54d/lib ] || export LD_LIBRARY_PATH="/work/build/CMSSW_12_4_X/el8_amd64_gcc10/tmp/el8_amd64_gcc10/external/tauolapp/1.1.8-5c4faaa3d5ddc2d72dd03a1f4952c54d/lib${LD_LIBRARY_PATH:+:$LD_LIBRARY_PATH}";

