# CMSSW_10_6_X/slc7_amd64_gcc700/tools/vincia/2.2.04-a7d48ce347d4f85ad6de04efb18b1fd9/etc/profile.d/dependencies-setup.csh
test $?GCC_ROOT != 0 || source /work/build/CMSSW_10_6_X/slc7_amd64_gcc700/tmp/slc7_amd64_gcc700/external/gcc/7.0.0-pafccj/etc/profile.d/init.csh
test $?ZLIB_X86_64_ROOT != 0 || source /work/build/CMSSW_10_6_X/slc7_amd64_gcc700/tmp/slc7_amd64_gcc700/external/zlib-x86_64/1.2.11-pafccj/etc/profile.d/init.csh
test $?OPENSSL_ROOT != 0 || source /work/build/CMSSW_10_6_X/slc7_amd64_gcc700/tmp/slc7_amd64_gcc700/external/openssl/1.0.2d-pafccj/etc/profile.d/init.csh
test $?ZLIB_ROOT != 0 || source /work/build/CMSSW_10_6_X/slc7_amd64_gcc700/tmp/slc7_amd64_gcc700/external/zlib-x86_64/1.2.11-pafccj/etc/profile.d/init.csh
test $?BZ2LIB_ROOT != 0 || source /work/build/CMSSW_10_6_X/slc7_amd64_gcc700/tmp/slc7_amd64_gcc700/external/bz2lib/1.0.6-pafccj/etc/profile.d/init.csh
test $?EXPAT_ROOT != 0 || source /work/build/CMSSW_10_6_X/slc7_amd64_gcc700/tmp/slc7_amd64_gcc700/external/expat/2.1.0-pafccj/etc/profile.d/init.csh
test $?DB6_ROOT != 0 || source /work/build/CMSSW_10_6_X/slc7_amd64_gcc700/tmp/slc7_amd64_gcc700/external/db6/6.0.30-pafccj/etc/profile.d/init.csh
test $?GDBM_ROOT != 0 || source /work/build/CMSSW_10_6_X/slc7_amd64_gcc700/tmp/slc7_amd64_gcc700/external/gdbm/1.10-pafccj/etc/profile.d/init.csh
test $?LIBFFI_ROOT != 0 || source /work/build/CMSSW_10_6_X/slc7_amd64_gcc700/tmp/slc7_amd64_gcc700/external/libffi/3.2.1-pafccj/etc/profile.d/init.csh
test $?SQLITE_ROOT != 0 || source /work/build/CMSSW_10_6_X/slc7_amd64_gcc700/tmp/slc7_amd64_gcc700/external/sqlite/3.22.0-pafccj/etc/profile.d/init.csh
test $?PYTHON_ROOT != 0 || source /work/build/CMSSW_10_6_X/slc7_amd64_gcc700/tmp/slc7_amd64_gcc700/external/python/2.7.14-pafccj/etc/profile.d/init.csh
test $?AUTOTOOLS_ROOT != 0 || source /work/build/CMSSW_10_6_X/slc7_amd64_gcc700/tmp/slc7_amd64_gcc700/external/autotools/1.3-pafccj/etc/profile.d/init.csh
test $?HEPMC_ROOT != 0 || source /work/build/CMSSW_10_6_X/slc7_amd64_gcc700/tmp/slc7_amd64_gcc700/external/hepmc/2.06.07-pafccj/etc/profile.d/init.csh
test $?LHAPDF_ROOT != 0 || source /work/build/CMSSW_10_6_X/slc7_amd64_gcc700/tmp/slc7_amd64_gcc700/external/lhapdf/6.2.3-a2a84f5990d32c24c7240b02577bf55e/etc/profile.d/init.csh
test $?PYTHIA8_ROOT != 0 || source /work/build/CMSSW_10_6_X/slc7_amd64_gcc700/tmp/slc7_amd64_gcc700/external/pythia8/240-97646d8095752c6e96f84bb0257efff7/etc/profile.d/init.csh
