# CMSSW_10_6_X/slc7_amd64_gcc700/tools/vincia/2.2.04-a7d48ce347d4f85ad6de04efb18b1fd9/etc/profile.d/init.sh
if [ -f /work/build/CMSSW_10_6_X/slc7_amd64_gcc700/tmp/slc7_amd64_gcc700/external/vincia/2.2.04-a7d48ce347d4f85ad6de04efb18b1fd9/etc/profile.d/dependencies-setup.sh ]; then . /work/build/CMSSW_10_6_X/slc7_amd64_gcc700/tmp/slc7_amd64_gcc700/external/vincia/2.2.04-a7d48ce347d4f85ad6de04efb18b1fd9/etc/profile.d/dependencies-setup.sh; fi
VINCIA_ROOT="/work/build/CMSSW_10_6_X/slc7_amd64_gcc700/tmp/slc7_amd64_gcc700/external/vincia/2.2.04-a7d48ce347d4f85ad6de04efb18b1fd9"
VINCIA_VERSION="2.2.04-a7d48ce347d4f85ad6de04efb18b1fd9"
VINCIA_REVISION="1"
VINCIA_CATEGORY="external"
[ ! -d /work/build/CMSSW_10_6_X/slc7_amd64_gcc700/tmp/slc7_amd64_gcc700/external/vincia/2.2.04-a7d48ce347d4f85ad6de04efb18b1fd9/bin ] || export PATH="/work/build/CMSSW_10_6_X/slc7_amd64_gcc700/tmp/slc7_amd64_gcc700/external/vincia/2.2.04-a7d48ce347d4f85ad6de04efb18b1fd9/bin${PATH:+:$PATH}";
[ ! -d /work/build/CMSSW_10_6_X/slc7_amd64_gcc700/tmp/slc7_amd64_gcc700/external/vincia/2.2.04-a7d48ce347d4f85ad6de04efb18b1fd9/lib ] || export LD_LIBRARY_PATH="/work/build/CMSSW_10_6_X/slc7_amd64_gcc700/tmp/slc7_amd64_gcc700/external/vincia/2.2.04-a7d48ce347d4f85ad6de04efb18b1fd9/lib${LD_LIBRARY_PATH:+:$LD_LIBRARY_PATH}";

