# CMSSW_12_4_15/el8_amd64_gcc10/tools/tauolapp/1.1.8-c2dff32d490b71605acef6b3c7cb0a36/etc/profile.d/init.sh
if [ -f /uscms_data/d3/pedrok/semivisible/production/forRun3/pythia/tmp/tmp_O9c7tXKHHr/el8_amd64_gcc10/external/tauolapp/1.1.8-c2dff32d490b71605acef6b3c7cb0a36/etc/profile.d/dependencies-setup.sh ]; then . /uscms_data/d3/pedrok/semivisible/production/forRun3/pythia/tmp/tmp_O9c7tXKHHr/el8_amd64_gcc10/external/tauolapp/1.1.8-c2dff32d490b71605acef6b3c7cb0a36/etc/profile.d/dependencies-setup.sh; fi
TAUOLAPP_ROOT="/uscms_data/d3/pedrok/semivisible/production/forRun3/pythia/tmp/tmp_O9c7tXKHHr/el8_amd64_gcc10/external/tauolapp/1.1.8-c2dff32d490b71605acef6b3c7cb0a36"
TAUOLAPP_VERSION="1.1.8-c2dff32d490b71605acef6b3c7cb0a36"
TAUOLAPP_REVISION="1"
TAUOLAPP_CATEGORY="external"
[ ! -d /uscms_data/d3/pedrok/semivisible/production/forRun3/pythia/tmp/tmp_O9c7tXKHHr/el8_amd64_gcc10/external/tauolapp/1.1.8-c2dff32d490b71605acef6b3c7cb0a36/bin ] || export PATH="/uscms_data/d3/pedrok/semivisible/production/forRun3/pythia/tmp/tmp_O9c7tXKHHr/el8_amd64_gcc10/external/tauolapp/1.1.8-c2dff32d490b71605acef6b3c7cb0a36/bin${PATH:+:$PATH}";
[ ! -d /uscms_data/d3/pedrok/semivisible/production/forRun3/pythia/tmp/tmp_O9c7tXKHHr/el8_amd64_gcc10/external/tauolapp/1.1.8-c2dff32d490b71605acef6b3c7cb0a36/lib ] || export LD_LIBRARY_PATH="/uscms_data/d3/pedrok/semivisible/production/forRun3/pythia/tmp/tmp_O9c7tXKHHr/el8_amd64_gcc10/external/tauolapp/1.1.8-c2dff32d490b71605acef6b3c7cb0a36/lib${LD_LIBRARY_PATH:+:$LD_LIBRARY_PATH}";

