# CMSSW_10_6_X/slc7_amd64_gcc700/tools/geneva/1.0-RC3-7d48a5948fe7fe1986bbecd35963d730/etc/profile.d/init.sh
if [ -f /work/build/CMSSW_10_6_X/slc7_amd64_gcc700/tmp/slc7_amd64_gcc700/external/geneva/1.0-RC3-7d48a5948fe7fe1986bbecd35963d730/etc/profile.d/dependencies-setup.sh ]; then . /work/build/CMSSW_10_6_X/slc7_amd64_gcc700/tmp/slc7_amd64_gcc700/external/geneva/1.0-RC3-7d48a5948fe7fe1986bbecd35963d730/etc/profile.d/dependencies-setup.sh; fi
GENEVA_ROOT="/work/build/CMSSW_10_6_X/slc7_amd64_gcc700/tmp/slc7_amd64_gcc700/external/geneva/1.0-RC3-7d48a5948fe7fe1986bbecd35963d730"
GENEVA_VERSION="1.0-RC3-7d48a5948fe7fe1986bbecd35963d730"
GENEVA_REVISION="1"
GENEVA_CATEGORY="external"
[ ! -d /work/build/CMSSW_10_6_X/slc7_amd64_gcc700/tmp/slc7_amd64_gcc700/external/geneva/1.0-RC3-7d48a5948fe7fe1986bbecd35963d730/bin ] || export PATH="/work/build/CMSSW_10_6_X/slc7_amd64_gcc700/tmp/slc7_amd64_gcc700/external/geneva/1.0-RC3-7d48a5948fe7fe1986bbecd35963d730/bin${PATH:+:$PATH}";
[ ! -d /work/build/CMSSW_10_6_X/slc7_amd64_gcc700/tmp/slc7_amd64_gcc700/external/geneva/1.0-RC3-7d48a5948fe7fe1986bbecd35963d730/lib ] || export LD_LIBRARY_PATH="/work/build/CMSSW_10_6_X/slc7_amd64_gcc700/tmp/slc7_amd64_gcc700/external/geneva/1.0-RC3-7d48a5948fe7fe1986bbecd35963d730/lib${LD_LIBRARY_PATH:+:$LD_LIBRARY_PATH}";
[ ! -d /work/build/CMSSW_10_6_X/slc7_amd64_gcc700/tmp/slc7_amd64_gcc700/external/geneva/1.0-RC3-7d48a5948fe7fe1986bbecd35963d730/${PYTHON_LIB_SITE_PACKAGES} ] || export PYTHON27PATH="/work/build/CMSSW_10_6_X/slc7_amd64_gcc700/tmp/slc7_amd64_gcc700/external/geneva/1.0-RC3-7d48a5948fe7fe1986bbecd35963d730/${PYTHON_LIB_SITE_PACKAGES}${PYTHON27PATH:+:$PYTHON27PATH}";

