# CMSSW_12_4_X/slc7_amd64_gcc10/tools/pythia8/306-2a52225c5b09dea44fcd9f82e1688682/etc/profile.d/dependencies-setup.csh
test $?GCC_ROOT != 0 || source /work/build/CMSSW_12_4_X/slc7_amd64_gcc10/tmp/slc7_amd64_gcc10/external/gcc/10.3.0-84898dea653199466402e67d73657f10/etc/profile.d/init.csh
test $?ZLIB_ROOT != 0 || source /work/build/CMSSW_12_4_X/slc7_amd64_gcc10/tmp/slc7_amd64_gcc10/external/zlib/1.2.11-71514f01d2850dfd0bcd391557259a58/etc/profile.d/init.csh
test $?BZ2LIB_ROOT != 0 || source /work/build/CMSSW_12_4_X/slc7_amd64_gcc10/tmp/slc7_amd64_gcc10/external/bz2lib/1.0.6-86270889250d9cb8193d1963706a39cb/etc/profile.d/init.csh
test $?EXPAT_ROOT != 0 || source /work/build/CMSSW_12_4_X/slc7_amd64_gcc10/tmp/slc7_amd64_gcc10/external/expat/2.1.0-e2e74851cebffc0aa0988d120ae74fcd/etc/profile.d/init.csh
test $?DB6_ROOT != 0 || source /work/build/CMSSW_12_4_X/slc7_amd64_gcc10/tmp/slc7_amd64_gcc10/external/db6/6.2.32-0f34c5862bf3098cbedece4ce9240238/etc/profile.d/init.csh
test $?GDBM_ROOT != 0 || source /work/build/CMSSW_12_4_X/slc7_amd64_gcc10/tmp/slc7_amd64_gcc10/external/gdbm/1.10-d603980cbd3ec2df3505d5b82740ee09/etc/profile.d/init.csh
test $?LIBUUID_ROOT != 0 || source /work/build/CMSSW_12_4_X/slc7_amd64_gcc10/tmp/slc7_amd64_gcc10/external/libuuid/2.34-fe91b21b8e54a48f2d45894d9da65f14/etc/profile.d/init.csh
test $?SQLITE_ROOT != 0 || source /work/build/CMSSW_12_4_X/slc7_amd64_gcc10/tmp/slc7_amd64_gcc10/external/sqlite/3.36.0-cdaed8943a1c7d35c738a99ec764027b/etc/profile.d/init.csh
test $?LIBFFI_ROOT != 0 || source /work/build/CMSSW_12_4_X/slc7_amd64_gcc10/tmp/slc7_amd64_gcc10/external/libffi/3.4.2-482b99521caf507051a7a5dce0eb5926/etc/profile.d/init.csh
test $?XZ_ROOT != 0 || source /work/build/CMSSW_12_4_X/slc7_amd64_gcc10/tmp/slc7_amd64_gcc10/external/xz/5.2.5-d6fed2038c4e8d6e04531d1adba59f37/etc/profile.d/init.csh
test $?PYTHON3_ROOT != 0 || source /work/build/CMSSW_12_4_X/slc7_amd64_gcc10/tmp/slc7_amd64_gcc10/external/python3/3.9.6-67e5cf5b4952101922f1d4c8474baa39/etc/profile.d/init.csh
test $?HEPMC_ROOT != 0 || source /work/build/CMSSW_12_4_X/slc7_amd64_gcc10/tmp/slc7_amd64_gcc10/external/hepmc/2.06.10-46867a6dcc6e5712b7953fe57085fcbd/etc/profile.d/init.csh
test $?HEPMC3_ROOT != 0 || source /work/build/CMSSW_12_4_X/slc7_amd64_gcc10/tmp/slc7_amd64_gcc10/external/hepmc3/3.2.5-c3cd50aeecf06b194814f1a75bf7872e/etc/profile.d/init.csh
test $?LHAPDF_ROOT != 0 || source /work/build/CMSSW_12_4_X/slc7_amd64_gcc10/tmp/slc7_amd64_gcc10/external/lhapdf/6.4.0-68defff11ffd434c73727d03802bfb85/etc/profile.d/init.csh
