# CMSSW_13_0_X/slc7_amd64_gcc11/tools/CICADA/1.2.2-eb162db096d0f432e5b2e16b8d5ecbea/etc/profile.d/init.sh
if [ -f /work/build/CMSSW_13_0_X/slc7_amd64_gcc11/tmp/slc7_amd64_gcc11/external/CICADA/1.2.2-eb162db096d0f432e5b2e16b8d5ecbea/etc/profile.d/dependencies-setup.sh ]; then . /work/build/CMSSW_13_0_X/slc7_amd64_gcc11/tmp/slc7_amd64_gcc11/external/CICADA/1.2.2-eb162db096d0f432e5b2e16b8d5ecbea/etc/profile.d/dependencies-setup.sh; fi
CICADA_ROOT="/work/build/CMSSW_13_0_X/slc7_amd64_gcc11/tmp/slc7_amd64_gcc11/external/CICADA/1.2.2-eb162db096d0f432e5b2e16b8d5ecbea"
CICADA_VERSION="1.2.2-eb162db096d0f432e5b2e16b8d5ecbea"
CICADA_REVISION="1"
CICADA_CATEGORY="external"
[ ! -d /work/build/CMSSW_13_0_X/slc7_amd64_gcc11/tmp/slc7_amd64_gcc11/external/CICADA/1.2.2-eb162db096d0f432e5b2e16b8d5ecbea/bin ] || export PATH="/work/build/CMSSW_13_0_X/slc7_amd64_gcc11/tmp/slc7_amd64_gcc11/external/CICADA/1.2.2-eb162db096d0f432e5b2e16b8d5ecbea/bin${PATH:+:$PATH}";
[ ! -d /work/build/CMSSW_13_0_X/slc7_amd64_gcc11/tmp/slc7_amd64_gcc11/external/CICADA/1.2.2-eb162db096d0f432e5b2e16b8d5ecbea/lib ] || export LD_LIBRARY_PATH="/work/build/CMSSW_13_0_X/slc7_amd64_gcc11/tmp/slc7_amd64_gcc11/external/CICADA/1.2.2-eb162db096d0f432e5b2e16b8d5ecbea/lib${LD_LIBRARY_PATH:+:$LD_LIBRARY_PATH}";

