# CMSSW_13_0_X/el8_amd64_gcc11/tools/hls/2019.08-bcaeca57fe2d934c73d5ebb7c1864fb8/etc/profile.d/dependencies-setup.csh
test $?GCC_ROOT != 0 || source /work/build/CMSSW_13_0_X/el8_amd64_gcc11/tmp/el8_amd64_gcc11/external/gcc/11.2.1-f9b9dfdd886f71cd63f5538223d8f161/etc/profile.d/init.csh
test $?GMAKE_ROOT != 0 || source /work/build/CMSSW_13_0_X/el8_amd64_gcc11/tmp/el8_amd64_gcc11/external/gmake/4.3-6228c6b5a92af64d32a9672586ea63b7/etc/profile.d/init.csh
