# CMSSW_13_0_X/el8_amd64_gcc11/tools/hls/2019.08-bcaeca57fe2d934c73d5ebb7c1864fb8/etc/profile.d/init.sh
if [ -f /work/build/CMSSW_13_0_X/el8_amd64_gcc11/tmp/el8_amd64_gcc11/external/hls/2019.08-bcaeca57fe2d934c73d5ebb7c1864fb8/etc/profile.d/dependencies-setup.sh ]; then . /work/build/CMSSW_13_0_X/el8_amd64_gcc11/tmp/el8_amd64_gcc11/external/hls/2019.08-bcaeca57fe2d934c73d5ebb7c1864fb8/etc/profile.d/dependencies-setup.sh; fi
HLS_ROOT="/work/build/CMSSW_13_0_X/el8_amd64_gcc11/tmp/el8_amd64_gcc11/external/hls/2019.08-bcaeca57fe2d934c73d5ebb7c1864fb8"
HLS_VERSION="2019.08-bcaeca57fe2d934c73d5ebb7c1864fb8"
HLS_REVISION="1"
HLS_CATEGORY="external"
[ ! -d /work/build/CMSSW_13_0_X/el8_amd64_gcc11/tmp/el8_amd64_gcc11/external/hls/2019.08-bcaeca57fe2d934c73d5ebb7c1864fb8/bin ] || export PATH="/work/build/CMSSW_13_0_X/el8_amd64_gcc11/tmp/el8_amd64_gcc11/external/hls/2019.08-bcaeca57fe2d934c73d5ebb7c1864fb8/bin${PATH:+:$PATH}";
[ ! -d /work/build/CMSSW_13_0_X/el8_amd64_gcc11/tmp/el8_amd64_gcc11/external/hls/2019.08-bcaeca57fe2d934c73d5ebb7c1864fb8/lib ] || export LD_LIBRARY_PATH="/work/build/CMSSW_13_0_X/el8_amd64_gcc11/tmp/el8_amd64_gcc11/external/hls/2019.08-bcaeca57fe2d934c73d5ebb7c1864fb8/lib${LD_LIBRARY_PATH:+:$LD_LIBRARY_PATH}";

