# CMSSW_12_4_15/slc7_amd64_gcc10/tools/pythia8/306-2a52225c5b09dea44fcd9f82e1688682/etc/profile.d/init.sh
if [ -f /uscms_data/d3/pedrok/semivisible/production/forRun3/pythia/tmp2/build/CMSSW_12_4_15/slc7_amd64_gcc10/tmp/slc7_amd64_gcc10/external/pythia8/306-2a52225c5b09dea44fcd9f82e1688682/etc/profile.d/dependencies-setup.sh ]; then . /uscms_data/d3/pedrok/semivisible/production/forRun3/pythia/tmp2/build/CMSSW_12_4_15/slc7_amd64_gcc10/tmp/slc7_amd64_gcc10/external/pythia8/306-2a52225c5b09dea44fcd9f82e1688682/etc/profile.d/dependencies-setup.sh; fi
PYTHIA8_ROOT="/uscms_data/d3/pedrok/semivisible/production/forRun3/pythia/tmp2/build/CMSSW_12_4_15/slc7_amd64_gcc10/tmp/slc7_amd64_gcc10/external/pythia8/306-2a52225c5b09dea44fcd9f82e1688682"
PYTHIA8_VERSION="306-2a52225c5b09dea44fcd9f82e1688682"
PYTHIA8_REVISION="1"
PYTHIA8_CATEGORY="external"
[ ! -d /uscms_data/d3/pedrok/semivisible/production/forRun3/pythia/tmp2/build/CMSSW_12_4_15/slc7_amd64_gcc10/tmp/slc7_amd64_gcc10/external/pythia8/306-2a52225c5b09dea44fcd9f82e1688682/bin ] || export PATH="/uscms_data/d3/pedrok/semivisible/production/forRun3/pythia/tmp2/build/CMSSW_12_4_15/slc7_amd64_gcc10/tmp/slc7_amd64_gcc10/external/pythia8/306-2a52225c5b09dea44fcd9f82e1688682/bin${PATH:+:$PATH}";
[ ! -d /uscms_data/d3/pedrok/semivisible/production/forRun3/pythia/tmp2/build/CMSSW_12_4_15/slc7_amd64_gcc10/tmp/slc7_amd64_gcc10/external/pythia8/306-2a52225c5b09dea44fcd9f82e1688682/lib ] || export LD_LIBRARY_PATH="/uscms_data/d3/pedrok/semivisible/production/forRun3/pythia/tmp2/build/CMSSW_12_4_15/slc7_amd64_gcc10/tmp/slc7_amd64_gcc10/external/pythia8/306-2a52225c5b09dea44fcd9f82e1688682/lib${LD_LIBRARY_PATH:+:$LD_LIBRARY_PATH}";

