# CMSSW_12_4_15/slc7_amd64_gcc10/tools/evtgen/2.0.0-46702aceecb30229742da6c34bfb6515/etc/profile.d/init.sh
if [ -f /uscms_data/d3/pedrok/semivisible/production/forRun3/pythia/tmp/tmp_UoATUWTxAG/slc7_amd64_gcc10/external/evtgen/2.0.0-46702aceecb30229742da6c34bfb6515/etc/profile.d/dependencies-setup.sh ]; then . /uscms_data/d3/pedrok/semivisible/production/forRun3/pythia/tmp/tmp_UoATUWTxAG/slc7_amd64_gcc10/external/evtgen/2.0.0-46702aceecb30229742da6c34bfb6515/etc/profile.d/dependencies-setup.sh; fi
EVTGEN_ROOT="/uscms_data/d3/pedrok/semivisible/production/forRun3/pythia/tmp/tmp_UoATUWTxAG/slc7_amd64_gcc10/external/evtgen/2.0.0-46702aceecb30229742da6c34bfb6515"
EVTGEN_VERSION="2.0.0-46702aceecb30229742da6c34bfb6515"
EVTGEN_REVISION="1"
EVTGEN_CATEGORY="external"
[ ! -d /uscms_data/d3/pedrok/semivisible/production/forRun3/pythia/tmp/tmp_UoATUWTxAG/slc7_amd64_gcc10/external/evtgen/2.0.0-46702aceecb30229742da6c34bfb6515/bin ] || export PATH="/uscms_data/d3/pedrok/semivisible/production/forRun3/pythia/tmp/tmp_UoATUWTxAG/slc7_amd64_gcc10/external/evtgen/2.0.0-46702aceecb30229742da6c34bfb6515/bin${PATH:+:$PATH}";
[ ! -d /uscms_data/d3/pedrok/semivisible/production/forRun3/pythia/tmp/tmp_UoATUWTxAG/slc7_amd64_gcc10/external/evtgen/2.0.0-46702aceecb30229742da6c34bfb6515/lib ] || export LD_LIBRARY_PATH="/uscms_data/d3/pedrok/semivisible/production/forRun3/pythia/tmp/tmp_UoATUWTxAG/slc7_amd64_gcc10/external/evtgen/2.0.0-46702aceecb30229742da6c34bfb6515/lib${LD_LIBRARY_PATH:+:$LD_LIBRARY_PATH}";

