# CMSSW_12_4_X/el8_amd64_gcc10/tools/evtgen/2.0.0-55fbd35f3034b8518f27ef588686f9a3/etc/profile.d/init.sh
if [ -f /work/build/CMSSW_12_4_X/el8_amd64_gcc10/tmp/el8_amd64_gcc10/external/evtgen/2.0.0-55fbd35f3034b8518f27ef588686f9a3/etc/profile.d/dependencies-setup.sh ]; then . /work/build/CMSSW_12_4_X/el8_amd64_gcc10/tmp/el8_amd64_gcc10/external/evtgen/2.0.0-55fbd35f3034b8518f27ef588686f9a3/etc/profile.d/dependencies-setup.sh; fi
EVTGEN_ROOT="/work/build/CMSSW_12_4_X/el8_amd64_gcc10/tmp/el8_amd64_gcc10/external/evtgen/2.0.0-55fbd35f3034b8518f27ef588686f9a3"
EVTGEN_VERSION="2.0.0-55fbd35f3034b8518f27ef588686f9a3"
EVTGEN_REVISION="1"
EVTGEN_CATEGORY="external"
[ ! -d /work/build/CMSSW_12_4_X/el8_amd64_gcc10/tmp/el8_amd64_gcc10/external/evtgen/2.0.0-55fbd35f3034b8518f27ef588686f9a3/bin ] || export PATH="/work/build/CMSSW_12_4_X/el8_amd64_gcc10/tmp/el8_amd64_gcc10/external/evtgen/2.0.0-55fbd35f3034b8518f27ef588686f9a3/bin${PATH:+:$PATH}";
[ ! -d /work/build/CMSSW_12_4_X/el8_amd64_gcc10/tmp/el8_amd64_gcc10/external/evtgen/2.0.0-55fbd35f3034b8518f27ef588686f9a3/lib ] || export LD_LIBRARY_PATH="/work/build/CMSSW_12_4_X/el8_amd64_gcc10/tmp/el8_amd64_gcc10/external/evtgen/2.0.0-55fbd35f3034b8518f27ef588686f9a3/lib${LD_LIBRARY_PATH:+:$LD_LIBRARY_PATH}";

