# CMSSW_12_4_15/slc7_amd64_gcc10/tools/pythia8/306-3535c0959cf5f0d91ee1acb501d730ff/etc/profile.d/init.sh
if [ -f /uscms_data/d3/pedrok/semivisible/production/forRun3/pythia/tmp/tmp_UoATUWTxAG/slc7_amd64_gcc10/external/pythia8/306-3535c0959cf5f0d91ee1acb501d730ff/etc/profile.d/dependencies-setup.sh ]; then . /uscms_data/d3/pedrok/semivisible/production/forRun3/pythia/tmp/tmp_UoATUWTxAG/slc7_amd64_gcc10/external/pythia8/306-3535c0959cf5f0d91ee1acb501d730ff/etc/profile.d/dependencies-setup.sh; fi
PYTHIA8_ROOT="/uscms_data/d3/pedrok/semivisible/production/forRun3/pythia/tmp/tmp_UoATUWTxAG/slc7_amd64_gcc10/external/pythia8/306-3535c0959cf5f0d91ee1acb501d730ff"
PYTHIA8_VERSION="306-3535c0959cf5f0d91ee1acb501d730ff"
PYTHIA8_REVISION="1"
PYTHIA8_CATEGORY="external"
[ ! -d /uscms_data/d3/pedrok/semivisible/production/forRun3/pythia/tmp/tmp_UoATUWTxAG/slc7_amd64_gcc10/external/pythia8/306-3535c0959cf5f0d91ee1acb501d730ff/bin ] || export PATH="/uscms_data/d3/pedrok/semivisible/production/forRun3/pythia/tmp/tmp_UoATUWTxAG/slc7_amd64_gcc10/external/pythia8/306-3535c0959cf5f0d91ee1acb501d730ff/bin${PATH:+:$PATH}";
[ ! -d /uscms_data/d3/pedrok/semivisible/production/forRun3/pythia/tmp/tmp_UoATUWTxAG/slc7_amd64_gcc10/external/pythia8/306-3535c0959cf5f0d91ee1acb501d730ff/lib ] || export LD_LIBRARY_PATH="/uscms_data/d3/pedrok/semivisible/production/forRun3/pythia/tmp/tmp_UoATUWTxAG/slc7_amd64_gcc10/external/pythia8/306-3535c0959cf5f0d91ee1acb501d730ff/lib${LD_LIBRARY_PATH:+:$LD_LIBRARY_PATH}";

